# oci-k8s-cluster/systemd/gdrive-sync.service
[Unit]
Description=Google Drive Backup Sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/sync_to_gdrive.sh
User=root
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7