# oci-k8s-cluster/systemd/k8s-health-check.service
[Unit]
Description=Kubernetes Cluster Health Watchdog
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/k8s-ops/cluster_health_check.sh --no-color
StandardOutput=append:/var/log/k8s-health-check.log
StandardError=append:/var/log/k8s-health-check.log
User=root
# 0=OK/warnings, 2=critical — do not treat warnings as unit failure
SuccessExitStatus=0 2
